# etc/profile.d/opt_julia.sh
#!/bin/bash

_user="$(id -u -n)"

if [ "$_user" == "researcher" ]; then
  export JULIA_PKGDIR=/opt/julia
fi
